# autostart/sequoia-rag.service
[Unit]
Description=Sequoia RAG Assistant
After=network.target sound.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/sorva_ToyotaRag
ExecStartPre=/bin/bash -c 'ollama serve > /tmp/ollama.log 2>&1 &'
ExecStartPre=/bin/sleep 3
ExecStart=/home/pi/sorva_ToyotaRag/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
Environment=WAKE_WORD_ENABLED=true

[Install]
WantedBy=multi-user.target
